# cf-backup/files/etc/crontab
* * * * * root /usr/bin/cron-job.sh /usr/bin/run-rack.sh >> /var/log/cron.log 2>&1
